# ansible/roles/install_agent/files/rstats.service
[Unit]
Description=Startup script for Rstats
Requires=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/openbach/agent/rstats/rstats.py
ExecReload=/usr/bin/python3 /opt/openbach/agent/rstats/rstats_reload.py

[Install]
WantedBy=multi-user.target

